# etc/10-keyboard.rules
ACTION=="bind", SUBSYSTEM=="hid", RUN+="/usr/bin/su breno -c /home/breno/bin/udev/keyboard.sh %p"
