# etc/40-wacom.rules
ACTION=="bind", SUBSYSTEM=="hid", RUN+="/usr/bin/su breno -c /home/breno/bin/udev/wacom.sh %p"
