# etc/40-wacom.rules
ACTION=="bind", SUBSYSTEM=="hid", RUN+="/usr/bin/su breno -c /opt/bin/wacom.sh %p"
