# files/home/.config/systemd/user/gpg-agent.target
[Unit]
Description=Start gpg-agent on login
Requires=gpg-agent.service
After=gpg-agent.service

[Install]
WantedBy=default.target
